# system/etc/udev/rules.d/90-tso-off.rules
ACTION=="add", SUBSYSTEM=="net", NAME=="wlp*", RUN+="/usr/bin/ethtool -K $name tso off"
