# system/etc/systemd/system/i3device-setup.service
[Unit]
Description=Setup devices

[Service]
Type=oneshot
ExecStart=/home/bigo/.bin/i3device-setup

[Install]
WantedBy=multi-user.target
