# home/.config/systemd/user/gotify-dunst.service
[Unit]
Description=Gotify-Dunst
After=network.target

[Service]
WorkingDirectory=/home/ANT.AMAZON.COM/matbigoi/github/gotify-dunst
ExecStart=/usr/bin/env python main.py
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
