# home/.config/systemd/user/wm.target
[Unit]
Description=Window manager target
BindsTo=graphical-session.target
